# ansible/roles/slurm/templates/supervisor_worker_slurm.conf.j2
[program:munge]
user=munge
command=/usr/sbin/munged --foreground
stdout_logfile=/var/log/munge/munge.log
redirect_stderr=true
[program:slurmd]
command=slurmd -D
